# src/main/ansible/roles/incron/files/incron.conf
description "incrond"
start on local-filesystems
stop on shutdown

exec incrond

respawn


