# src/main/ansible/roles/box-cifs-store/files/mount-box.conf

start on (filesystem and started networking)
stop on shutdown

task
respawn

exec /sbin/mount-box.hs

post-stop exec sleep 5



