# src/main/ansible/roles/box-common/files/mount-box-home.conf
description "Mount box home"

start on filesystem
stop on shutdown

task

exec /sbin/mount-box-home.sh

